# tests/script-init/supervisord/resourced.conf
[program:resourced]
environment=RESOURCED_CONFIG_DIR="/go/src/github.com/resourced/resourced/tests/resourced-configs"
command=/go/bin/resourced
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
startsecs=2
startretries=3
user=root
stdout_logfile=/var/log/resourced-agent.stdout.log
stderr_logfile=/var/log/resourced-agent.stderr.log

; Using superlance memmon is recommended since memory can grow as you mine more server data.
[eventlistener:resourced]
command=memmon -p resourced=160MB
events=TICK_60