# tests/script-init/upstart/resourced.conf
# resourced - resourced job file

description "ResourceD Upstart script"
author "Didip Kerabat <didipk@gmail.com>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn
respawn limit 10 5

# Start the process
exec su root -c 'RESOURCED_CONFIG_DIR=/go/src/github.com/resourced/resourced/tests/resourced-configs /go/bin/resourced'